# launchd/com.mikkegt.ai-usage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mikkegt.ai-usage</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>$HOME/workspace/github.com/mikkegt/mikkegt/scripts/run_daily.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <!-- 03:00 JST = 18:00 UTC.
         launchd は LOCAL TIME 解釈なので JST マシン上では 03:00 を指定する。 -->
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/com.mikkegt.ai-usage.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/com.mikkegt.ai-usage.err.log</string>
</dict>
</plist>
